# src-tauri/scripts/archdeck.service.template
[Unit]
Description={{APP_NAME}} - Stream Deck Controller
Documentation=https://github.com/yourusername/{{APP_NAME_LOWER}}
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart={{INSTALL_PATH}}/{{APP_NAME_LOWER}} --hidden
Restart=on-failure
RestartSec=5

# Environment for GUI applications
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0

[Install]
WantedBy=default.target
